# deploy/sri-visor/sri-websockify.service
[Unit]
Description=websockify del visor SRI (127.0.0.1:6080 -> 127.0.0.1:5900)
After=sri-x11vnc.service
Requires=sri-x11vnc.service

[Service]
User=www-data
Group=www-data
# Puente WebSocket <-> VNC, solo en loopback. nginx lo expone con validación de token.
ExecStart=/usr/bin/websockify 127.0.0.1:6080 127.0.0.1:5900
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
