# deploy/sri-visor/sri-x11vnc.service
[Unit]
Description=x11vnc del display :99 (descarga asistida SRI)
After=sri-xvfb.service
Requires=sri-xvfb.service

[Service]
User=www-data
Group=www-data
Environment=DISPLAY=:99
# Esperar a que Xvfb :99 esté listo (necesita x11-utils -> xdpyinfo)
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do xdpyinfo -display :99 >/dev/null 2>&1 && exit 0; sleep 1; done; exit 1'
# Solo localhost; sin contraseña (el acceso lo controla el token validado por nginx).
ExecStart=/usr/bin/x11vnc -display :99 -forever -shared -localhost -rfbport 5900 -nopw -noxdamage -quiet
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
